# packaging/9.0.x/debian/valkey.service
[Service]
Type=simple
EnvironmentFile=-/etc/default/valkey
ExecStart=/usr/bin/valkey --modules-dir $VALKEY_MODULES_DIR
